# .docker/etc/supervisor.d/php-fpm.ini
[group:php-fpm]
programs=php-fpmd
priority=20

[program:php-fpmd]
process_name=%(program_name)s_%(process_num)02d
command=php-fpm -y %(ENV_DOCKER_ETC_DIR)s/php%(ENV_PHPV)s/php-fpm.conf
startsecs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0